# automation/systemd/bonsaipr-automation.service
[Unit]
Description=BonsaiPR Automation Service
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/weekly-bonsaipr-automation
ExecStart=/usr/bin/python3 /path/to/weekly-bonsaipr-automation/src/main.py
Restart=on-failure

[Install]
WantedBy=multi-user.target