# cron/weekly-automation.cron
# This file contains the cron job configuration to schedule the execution of the automation scripts weekly.

# Run the automation every Sunday at 2 AM
# Logs are written to the project's logs directory with timestamp
0 2 * * 0 cd /home/falken10vdl/bonsaiPRDevel/weekly-bonsaipr-automation && /usr/bin/python3 src/main.py >> logs/automation-$(date +\%Y\%m\%d-\%H\%M\%S).log 2>&1
